# debian/postrm
#!/bin/sh

set -e

if [ -f /usr/share/debconf/confmodule ]; then
    . /usr/share/debconf/confmodule
fi

if [ -f /usr/share/dbconfig-common/dpkg/postrm ]; then
    . /usr/share/dbconfig-common/dpkg/postrm
    dbc_go yubikey-val $@
fi

if [ "$1" = "purge" ]; then
    #This file is generated in postinst
    rm -f /etc/ykval/config-db.php
    if which ucf >/dev/null 2>&1; then
        ucf --purge /etc/ykval/config-db.php
    fi
fi

#if [ "$1" = "purge" ]; then
#    # these files are generated in postinst
#    for cf in /etc/ykksm/config-db.php /etc/ykksm/config-db.cfg; do
#	rm -f $cf
#	if which ucf >/dev/null 2>&1; then
#	    ucf --purge $cf
#	fi
#    done
#fi

#if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
#    a2dissite yubikey-ksm yubikey-ksm-ssl || true
#
#    if [ -x /usr/sbin/invoke-rc.d ]; then
#	invoke-rc.d apache2 reload 3>/dev/null || true
#    else
#	/etc/init.d/apache2 reload 3>/dev/null || true
#    fi
#fi

#DEBHELPER#

exit 0
